# cron/crontabs/reset_collections
# min   hour    day   month   weekday   command
*/10    *       *     *       *         wget -qO - --header "x-api-secret: ${API_SECRET}" http://api:3000/api/reset

